# MyModule/UniCron/UniCrond.cron
* * * * * /data/adb/modules/UniCron/UniCrond.sh
* * * * 1,3,5 rm -f /data/adb/modules/logs/UniCron.log