# data/adb/modules/Demo/UniCron/Demo.cron
* * * * * /mnt/d/HOME/ARCHLINNUX/GitHub/RootManage-Module-Model/data/adb/modules/Dome/Demo.sh