# MyModule/cron/crontabs/root
* * * * * /data/adb/modules/MagicNet/test.sh
* * * * * /data/adb/modules/UniCron/UniCron.sh
* 6 * * * /data/adb/modules/UniCron/logCleaner.sh
