# MyModule/API/cron_tasks/MagicNet.cron
* * * * * /data/adb/modules/MagicNet/test.sh