# MyModule/API/cron_tasks/UniCron.cron
* * * * * /data/adb/modules/UniCron/service